# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/run/supervisord.pid

# ── PostgreSQL ─────────────────────────────────────────────────
[program:postgresql]
command=/usr/lib/postgresql/16/bin/postgres
    -D /var/lib/postgresql/16/main
    -c config_file=/etc/postgresql/16/main/postgresql.conf
user=postgres
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

# ── Bun.js backend ─────────────────────────────────────────────
[program:backend]
command=/usr/local/bin/bun run /app/src/index.ts
directory=/app
autostart=true
autorestart=true
priority=20
startretries=5
startsecs=3
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=
    DATABASE_URL="postgresql://solum:%(ENV_DB_PASSWORD)s@127.0.0.1:5432/officemate",
    JWT_SECRET="%(ENV_JWT_SECRET)s",
    PORT="3000",
    NODE_ENV="production",
    ALLOWED_ORIGINS="http://localhost,http://127.0.0.1"

# ── Nginx frontend ─────────────────────────────────────────────
[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=30
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
